# src/jobs/.crontab
* * * * *  /usr/bin/curl node $WORKDIR/notifications/index.ts >> $WORKDIR/jobs.log 2>&1
